# sysproduction/linux/crontab
################## The big daily show #######################
# Use this for stuff to run on startup
#
# @reboot
#
#
# Daily sample FX Prices at 6 am
#
0 6  * * 1-5 $HOME/workspace3/psystemtrade/sysproduction/linux/p run.py sysproduction.updateFxPrices.update_fx_prices >> $HOME/echos/samplefxprices 2>&1
#
